# slapsync.timer
[Unit]
Description=Run Slapsync daily

[Timer]
OnCalendar=*-*-* 1:00:00
Persistent=true 

[Install]
WantedBy=timers.target
